# .circleci/kafka/supervisor/create-topics.conf
[program:create-topics]
command=/usr/bin/create-topics.sh
autostart=true
autorestart=false
